# src/config/autostart/terminal.desktop
# Start a terminal emulator on a GUI session launch

[Desktop Entry]
Type=Application
Name=Terminal Emilator
NoDisplay=true
Exec=lxterminal

